# deprecated_jcz_bridge/systemd/gadget-setup.service
[Unit]
Description=BJJCZ USB Gadget Setup (configfs + FunctionFS)
DefaultDependencies=no
After=sys-kernel-config.mount systemd-modules-load.service
Before=jcz-bridge.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /opt/d1ultra-bridge/setup_gadget.sh

[Install]
WantedBy=multi-user.target
