# jcz_bridge/systemd/jcz-bridge.service
[Unit]
Description=JCZ-to-D1 Ultra Bridge (BJJCZ emulator)
After=gadget-setup.service network.target
Requires=gadget-setup.service

[Service]
Type=simple
WorkingDirectory=/opt/d1ultra-bridge
ExecStart=/usr/bin/python3 /opt/d1ultra-bridge/jcz_bridge.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
